# package/debian/postinst
#!/bin/sh
set -e

echo "config file location: /etc/d2hs/d2hs.json"

systemctl daemon-reload

systemctl enable DNS2HostsSyncer.timer
systemctl start DNS2HostsSyncer.timer
