# systemd/DNS2HostsSyncer.service
[Unit]
Description=DNS2HostsSyncer

[Service]
Type=oneshot
ExecStart=intentionally leave blank

[Install]
WantedBy=multi-user.target
